# gandive-bot.service
[Unit]
Description=GandiveBot — Premium Crypto Signal Telegram Bot
Documentation=https://t.me/GandiveBot
After=network.target

[Service]
Type=simple
User=enishshah2
WorkingDirectory=/home/enishshah2/gandive_bot
ExecStart=/home/enishshah2/gandive_bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/enishshah2/gandive_bot/logs/bot.log
StandardError=append:/home/enishshah2/gandive_bot/logs/bot.error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/enishshah2/gandive_bot
ProtectHome=false

[Install]
WantedBy=multi-user.target
